# web15.conf
[program:web15]
command=/usr/local/bin/gunicorn wsgi --bind 0.0.0.0:80
directory=/root/web15
autostart=true